# syslog-ng.upstart

start on pld.sysinit-done
stop on pld.shutdown-started

env SERVICE_syslog=y
export SERVICE_syslog

respawn

console output

exec @@SBINDIR@@/syslog-ng --process-mode=background -f /etc/syslog-ng/syslog-ng.conf
expect fork

# vi: ft=upstart
